# puppet-orchestration/modules/flask_bootstrap/templates/uwsgi.conf.erb
[supervisord]

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini <%= scope.lookupvar("flask_bootstrap::app_directory") %>/uwsgi.ini

autostart=true
autorestart=true
#stderr_logfile = /tmp/uwsgi_err.log
#stdout_logfile = /tmp/uwsgi_out.log
stopsignal=INT
